# Systemd Timer/FileSync.timer
[Unit]
Description=FileSync

[Timer]
OnCalendar=00/1:50
Persistent=true

[Install]
WantedBy=timers.target